# deploy/dd_evn_deploy/nginx.cron
0 0 * * * root /usr/sbin/logrotate -f /usr/local/nginx/conf/nginx.logrotate
